# .etc/system/streamlit-dev.service
[Unit]
Description=Streamlit OCR App (dev)
After=network.target
 
[Service]
User=rb_admin2
WorkingDirectory=/home/rb_admin2/apps/main-dev/rb-ocr
# Optional isolation for outputs/artifacts if supported by app
Environment="PATH=/home/rb_admin2/.venv/streamlit/bin"
# Example optional: Environment=RUNS_ROOT=/home/rb_admin2/apps/main-dev/rb-ocr/runs
ExecStart=/home/rb_admin2/.venv/streamlit/bin/streamlit run /home/rb_admin2/apps/main-dev/rb-ocr/app.py --server.port 8006 --server.address localhost
Restart=always
 
[Install]
WantedBy=multi-user.target
